# config/etc/profile.d/fcitx.sh
# ======== ime
# Ref:
#     https://fcitx-im.org/wiki/Setup_Fcitx_5#.7E.2F.xprofile
#     https://wiki.archlinux.org/title/Fcitx5#Usage
#     https://wiki.archlinux.org/title/Fcitx#Usage
export XMODIFIERS=@im=fcitx
export GTK_IM_MODULE=fcitx
export QT_IM_MODULE=fcitx
export SDL_IM_MODULE=fcitx
export INPUT_METHOD=fcitx
